# meta-autonomos-raspberrypi/recipes-core/mount-fdir/files/mount-fdir.service
[Unit]
Description=Mount FDIR - Verify and retry critical mounts
Documentation=file:///etc/mount-fdir.conf
DefaultDependencies=no
After=local-fs.target
Before=sysinit.target

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/sbin/mount-fdir.sh

# Don't fail the boot if mounts fail - allow degraded operation
# Dependent services use ConditionPathIsMountPoint for graceful handling
SuccessExitStatus=0 1

[Install]
WantedBy=sysinit.target
