# meta-autonomos-core/recipes-core/read-only-rootfs-config/files/data-journal.service
[Unit]
Description=Bind mount persistent journal logs
DefaultDependencies=no
Before=systemd-journald.service
After=mount-fdir.service
ConditionPathIsMountPoint=/data

[Service]
Type=oneshot
RemainAfterExit=yes
# Create directories ourselves (can't depend on tmpfiles - creates ordering cycle with journald)
ExecStartPre=/bin/mkdir -p /data/log/journal
ExecStartPre=/bin/mkdir -p /var/volatile/log/journal
ExecStart=/bin/mount --bind /data/log/journal /var/volatile/log/journal

[Install]
WantedBy=sysinit.target
