# meta-autonomos-core/recipes-core/read-only-rootfs-config/files/data-ssh-keys.service
[Unit]
Description=Bind mount persistent SSH host keys
DefaultDependencies=no
Before=sshdgenkeys.service sshd.service sshd.socket
After=mount-fdir.service
ConditionPathIsMountPoint=/data

[Service]
Type=oneshot
RemainAfterExit=yes
# Create directory structure and bind mount
# Keys will be generated here by sshdgenkeys on first boot
ExecStartPre=/bin/mkdir -p /data/etc/ssh
ExecStartPre=/bin/mkdir -p /run/ssh
ExecStartPre=/bin/chmod 700 /data/etc/ssh
ExecStart=/bin/mount --bind /data/etc/ssh /run/ssh

[Install]
WantedBy=sysinit.target
